# .config/systemd/user/homebrew.cliproxyapi.service.d/override.conf
[Service]
ExecStart=
ExecStart="/home/linuxbrew/.linuxbrew/opt/cliproxyapi/bin/cliproxyapi" -config "%h/.config/cliproxyapi/cliproxyapi.conf"
